# container/files/etc/supervisor.d/app.ini
[program:node-app]
directory = /app/server
command = npm start
process_name = app-%(process_num)s
numprocs=1
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/var/log/app-out-%(process_num)s.log
stdout_logfile_maxbytes=10KB
stderr_logfile=/var/log/app-error-%(process_num)s.log
stderr_logfile_maxbytes=10KB
